# deploy/systemd/smartbling-backend.service
[Unit]
Description=smartBling FastAPI backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/smartB
Environment=BACKEND_HOST=0.0.0.0
Environment=BACKEND_PORT=8000
Environment=BACKEND_LOG_LEVEL=info
ExecStart=/usr/bin/env bash /opt/smartB/scripts/run-backend-prod.sh

# Auto-restart com backoff progressivo (evita thundering herd)
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Limita memória para evitar OOM (deixa ~1GB para sistema com 4GB total)
MemoryMax=3072M
MemoryAccounting=yes

# Kill gracioso: 30s para shutdown limpo, depois SIGKILL
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30


[Install]
WantedBy=multi-user.target
